# server/rpi_cam-supervisor.conf.tmpl
[program:rpi_cam]
numprocs = 1
numprocs_start = 1
process_name = rpi_cam

command={{ command }} {{ args }}

user=pi
autostart=true
autorestart=true